# ops/ansible/roles/spark-master/templates/sparkmaster.conf
description "Spark Master"

# automatically start
start on filesystem

# auto respawn service
respawn
#respawn limit 15 5

# start script TODO: replace by yarn
script
    export SPARK_MASTER_IP=$(curl -s http://169.254.169.254/latest/meta-data/public-hostname)
    printf "SPARK_MASTER_IP=$SPARK_MASTER_IP"
    exec {{spark_home}}/sbin/start-master.sh
end script